# mac-agent/com.scladpro.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.scladpro.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <!-- ЗАМЕНИ путь на свой -->
    <string>/Users/ИМЯ/scladpro/mac-agent/run-agent.sh</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/ИМЯ/scladpro/mac-agent/agent.log</string>
  <key>StandardErrorPath</key><string>/Users/ИМЯ/scladpro/mac-agent/agent.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
